# install/atles_remote.service
[Unit]
Description=ATLeS RPC service
After=sshd.service
Requires=sshd.service

[Service]
WorkingDirectory=/home/pi/ATLeS
ExecStart=/home/pi/ATLeS/src/atles_remote.py
Type=idle

[Install]
WantedBy=multi-user.target
